# window-icon-updater/qubes-icon-receiver.socket
[Socket]
ListenStream=/run/qubes/icon-receiver.sock
PartOf=qubes-session.target

[Install]
WantedBy=qubes-session.target
